# deploy/agent-gateway/systemd/agent-gateway@.service
[Unit]
Description=Agent gateway profile: %i
Documentation=file:/opt/agent/gateways/%i/env
After=network-online.target
Wants=network-online.target
# Refuse to start with no profile config or no runtime selected (fail loudly, never run stale):
AssertPathExists=/opt/agent/gateways/%i/env
AssertPathIsSymbolicLink=/opt/agent/gateways/%i/current
# NOTE: deliberately NO Requires=/Wants=/PartOf=/BindsTo= on other gateway profiles.
# Each profile starts, stops, and fails completely independently.

[Service]
Type=simple
# Production-isolation mode: each profile runs as its OWN service account — this (plus the hardening
# below) is the file-security boundary. App-state profiles under ONE Unix user are NOT a security
# boundary. For convenience/dev mode, use User=agent / Group=agent instead.
User=agent-%i
Group=agent-%i

# Optional shared overlays load FIRST; the per-profile env is authoritative and loads LAST
# (systemd: a later EnvironmentFile wins on conflict). The runtime-source vars live ONLY in the
# profile env, so this ordering keeps that file as the single source of truth.
EnvironmentFile=-/etc/agent/model.env
EnvironmentFile=/opt/agent/gateways/%i/env
WorkingDirectory=/opt/agent/gateways/%i/current

# Shared launcher: logs an alignment banner (pid/cwd/PYTHONPATH/runtime/sha), then exec()s the
# resolved entrypoint. It does NOT set the runtime — it only reads it from the env file.
ExecStart=/opt/agent/_lib/launch

Restart=on-failure
RestartSec=2
TimeoutStartSec=60
SyslogIdentifier=agent-gw-%i
StandardOutput=append:/opt/agent/gateways/%i/logs/gateway.log
StandardError=append:/opt/agent/gateways/%i/logs/gateway.err

# --- per-profile OS isolation (production mode): writes only its own home/ and logs/; the runtime
# tree (current -> immutable release) stays read-only under ProtectSystem=strict ---
UMask=0077
ReadWritePaths=/opt/agent/gateways/%i/home /opt/agent/gateways/%i/logs
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
RestrictSUIDSGID=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# Default-deny egress; allow loopback + the approved local proxy only (tune per profile).
IPAddressDeny=any
IPAddressAllow=127.0.0.1/32 ::1/128 10.0.0.1/32
# Resource limits (tune during canary).
TasksMax=512
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
